# deployment/ec2/callwaitingai.service
# Systemd service file for CallWaitingAI on EC2
[Unit]
Description=CallWaitingAI Docker Compose Application
Requires=docker.service
After=docker.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/ec2-user/callwaitingai
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
ExecReload=/usr/local/bin/docker-compose restart
TimeoutStartSec=300
TimeoutStopSec=120
User=ec2-user
Group=ec2-user

# Restart policy
Restart=on-failure
RestartSec=30
StartLimitBurst=3

# Environment
Environment=COMPOSE_HTTP_TIMEOUT=300
Environment=DOCKER_CLIENT_TIMEOUT=300

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/ec2-user/callwaitingai

[Install]
WantedBy=multi-user.target